# systemd/Service.service.in
[Unit]
Description=@PROJECT_DESCRIPTION@
Wants=systemd-udevd.service
After=systemd-udevd.service

[Service]
Type=dbus
BusName=@DBUS_BUS_NAME@
ExecStart=@SERVICE_EXECUTABLE@
